# template/mysqldump.crontab
# MIN HOUR DAY MONTH DAYOFWEEK COMMAND
10 00 * * * . /etc/profile; mysqldump_all.sh | mail -s "Appaserver Mysqldump Database" yourname@foo.com
